# debian/postinst
#!/bin/sh

set -e

# Register initscript
update-rc.d shutdown-kvm-guests.sh stop 20 016

# Continue with debhelper stuff
#DEBHELPER#
