# debian/postinst
#!/bin/sh

set -e

# Register initscript
update-rc.d shutdown-kvm-guests.sh stop 0 1 6

# Continue with debhelper stuff
#DEBHELPER#
