# debian/shutdown-kvm-guest/DEBIAN/postinst
#!/bin/sh

set -e

# Register initscript
update-rc.d shutdown-kvm-guests.sh stop 06

# Continue with debhelper stuff

