# debian/shutdown-kvm-guest/DEBIAN/postrm
#!/bin/sh

set -e

# Register initscript
update-rc.d shutdown-kvm-guests.sh remove

# Continue with debhelper stuff

